# delhi-electricity/crontab
*/1 * * * * (cd / && /usr/local/bin/python3 -m parser.IN_DL)
