# delhi-electricity/crontab
*/1 * * * * (cd / && python3 -m parser.IN_DL)
